# roles/mail/files/msmtp-queue.service
[Unit]
Description=Send queued mails

[Service]
Type=oneshot
ExecStart=/usr/local/bin/msmtp-queue -r
