# roles/mail/files/msmtp-queue.timer
[Unit]
Description=Send queue mails

[Timer]
OnCalendar=*:0/15
AccuracySec=1m

[Install]
WantedBy=multi-user.target