# roles/docker/files/docker-prune.service
[Unit]
Description= Prune dangeling and unused docker objects

[Service]
Type=oneshot
ExecStart=docker system prune --filter 'until=168h' --force