# roles/restic/files/restic-prune@.service
[Unit]
Description= Prune restic repo from %i

[Service]
Type=oneshot
EnvironmentFile=/etc/restic/%i
ExecStart=restic prune