# roles/restic/files/restic-prune@.timer
[Unit]
Description= Trigger restic prune for %i

[Timer]
OnCalendar=weekly
AccuracySec=1h
Persistent=True

[Install]
WantedBy=multi-user.target