# roles/image_builder/files/image-builder-cleanup.timer
[Unit]
Description= Clean image-builder garbage files

[Timer]
OnCalendar=*-*-* 0:45:00

[Install]
WantedBy=timers.target