# roles/image_builder/files/image-builder@.timer
[Unit]
Description= Build custom docker image

[Timer]
OnCalendar=*-*-* 0:00:00

[Install]
WantedBy=timers.target